# daemon/packaging/macos/launchd/com.centy.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.centy.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/centy-daemon</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CENTY_DAEMON_ADDR</key>
        <string>127.0.0.1:50051</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/centy-daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/centy-daemon.error.log</string>

    <key>WorkingDirectory</key>
    <string>/usr/local/var/lib/centy</string>
</dict>
</plist>
